# etc/supervisor/conf.d/celery-translate.conf
; supervisor config file

[unix_http_server]
file=/var/run/supervisor.sock   ; (the path to the socket file)
chmod=0700                       ; sockef file mode (default 0700)
; The only purpose of this file is to silent "CRITICAL" error
; when starting up, see https://github.com/Supervisor/supervisor/issues/694
username = dummy
password = dummy

[program:celery-translate]
environment = CELERY_WORKER_RUNNING=1,CELERY_APP=weblate.utils
command = /usr/local/bin/celery worker --hostname 'translate@%%h' --loglevel info --queues=translate --prefetch-multiplier=4 %(ENV_CELERY_TRANSLATE_OPTIONS)s
stdout_events_enabled=true
stderr_events_enabled=true
