# etc/supervisor/conf.d/celery-memory.conf
[program:celery-memory]
environment = CELERY_WORKER_RUNNING=1,CELERY_APP=weblate.utils
command = /usr/local/bin/celery worker --hostname 'memory@%%h' --loglevel info --queues=memory --prefetch-multiplier=10 %(ENV_CELERY_MEMORY_OPTIONS)s
stdout_events_enabled=true
stderr_events_enabled=true
