# etc/supervisor/conf.d/celery-beat.conf
[program:celery-beat]
environment = CELERY_WORKER_RUNNING=1,CELERY_APP=weblate.utils
command = /usr/local/bin/celery beat --loglevel info --pidfile /run/celery/beat.pid %(ENV_CELERY_BEAT_OPTIONS)s
stdout_events_enabled=true
stderr_events_enabled=true
